# etc/tasks.ini
[group:tasks]
programs=counter_clear,put_tasks_toutiao,put_tasks_weibo,put_tasks_weixin,sogou_cookies,weixin_cookies


[program:counter_clear]
command=python tasks/run_job_counter_clear.py
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=logs/counter_clear.log


[program:put_tasks_toutiao]
command=python tasks/run_job_put_tasks_toutiao.py
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=logs/put_tasks_toutiao.log


[program:put_tasks_weibo]
command=python tasks/run_job_put_tasks_weibo.py
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=logs/put_tasks_weibo.log


[program:put_tasks_weixin]
command=python tasks/run_job_put_tasks_weixin.py
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=logs/put_tasks_weixin.log


[program:sogou_cookies]
command=python tasks/run_job_sogou_cookies.py
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=logs/sogou_cookies.log


[program:weixin_cookies]
command=python tasks/run_job_weixin_cookies.py
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=logs/weixin_cookies.log
